# conf/supervisord.conf
[supervisord]
pidfile=/tmp/supervisor.pid
logfile=/var/log/supervisor.log
nodaemon=true
user=root

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:uwsgi]
command=uwsgi --ini %(ENV_PROJECT_DIR)s/%(ENV_UWSGI_INI_FILENAME)s
stdout_logfile=/var/log/uwsgi.log
stderr_logfile=/var/log/uwsgi_error.log
autostart=false
autorestart=true
priority=4

[program:nginx]
command=nginx %(ENV_NGINX_OPTIONS)s
stdout_logfile=/var/log/nginx.log
stderr_logfile=/var/log/nginx_error.log
autostart=false
autorestart=true
priority=5

[program:memcached]
command=memcached -m %(ENV_MEMCACHED_MEMORY)s -s %(ENV_MEMCACHED_SOCK)s -P %(ENV_MEMCACHED_PID)s -u %(ENV_MEMCACHED_GROUP)s -vv
stdout_logfile=/var/log/memcached.log
stderr_logfile=/var/log/memcached_error.log
autostart=false
autorestart=true
priority=3

[program:initialization]
autostart=false
autorestart=false
stdout_logfile=/var/log/initialization.log
stderr_logfile=/var/log/initialization_error.log
command=bash -x %(ENV_SOURCES_DIR)s/conf/%(ENV_INITIALIZATION_SCRIPT)s
priority=2

[program:startup]
autostart=true
autorestart=false
stdout_logfile=/var/log/startup.log
stderr_logfile=/var/log/startup_error.log
priority=1
command=bash -x %(ENV_SOURCES_DIR)s/conf/start.supervisord
